# infra/roles/node_exporter/templates/node-exporter.service
[Service]
Type=simple
ExecStart=/usr/bin/node_exporter
